# lib/systemd/system/nmosconnection.service
[Unit]
Description=nmosconnectionmanagement
After=network-online.target

[Service]
ExecStart=/usr/bin/connectionmanagement
