# reverse-proxy/lib/systemd/system/proxylisting.service
[Unit]
Description=proxylisting
After=network-online.target

[Service]
ExecStart=/usr/bin/proxylisting
